# configs/systemd/atom-manager.service
# {{HOSTNAME}} System Manager API — systemd service
#
# Install: sudo cp {{HOSTNAME}}-manager.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now {{HOSTNAME}}-manager

[Unit]
Description={{HOSTNAME}} System Manager API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/{{HOSTNAME}}-manager
ExecStart=/opt/{{HOSTNAME}}-manager/venv/bin/uvicorn app:app --host 127.0.0.1 --port 9000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
